# networktablet.service
[Unit]
Description=GfxTablet Service

[Service]
Type=simple
ExecStart=/usr/bin/networktablet
